# systemd/user/randomwall.service
# ~/.config/systemd/user/randomwall.service
[Unit]
Description=Change GNOME wallpaper

[Service]
ExecStart=/usr/local/bin/randomwall
