# Config/autostart/myterm.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/myterm
