# Config/autostart/terminal.desktop
[Desktop Entry]
Type=Application
Exec=qterminal --name qterminal --title QTerminal --geometry -20-75
