# Config/autostart/variety-profile.desktop
[Desktop Entry]
Type=Application
Exec=echo "/home/ronnie/.config/variety-fractals/" > /home/ronnie/.variety-profile
