# Config/autostart/xbindkeys.desktop
[Desktop Entry]
Type=Application
Exec=xbindkeys -f ~/.xbindkeysrc
