# Config/autostart/xcompmgr.desktop
[Desktop Entry]
Type=Application
Exec=xcompmgr -c
